# conf/crontab
*	*	*	*	*	www		/root/ss_hotspot.pl --action=update
